# roles/lunaria/files/lunaria-frames.service
[Unit]
Description=lunaria page renderer and rotator (brief -> frame.png)
After=network-online.target
Wants=network-online.target

[Service]
User=lunaria
WorkingDirectory=/var/lib/lunaria
EnvironmentFile=/etc/default/lunaria
ExecStart=/usr/local/bin/lunaria-frames
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
