# roles/lunaria/files/mediamtx.service
[Unit]
Description=mediamtx media server (lunaria HLS transport)
After=network-online.target
Wants=network-online.target

[Service]
User=lunaria
WorkingDirectory=/opt/mediamtx
ExecStart=/opt/mediamtx/mediamtx /opt/mediamtx/mediamtx.yml
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
