# roles/pub/templates/publish-morning-brief.timer.j2
[Unit]
Description=Run publish-morning-brief on the Morning Brief publish schedule

[Timer]
OnCalendar={{ pub_timer_oncalendar }}
Persistent=true
RandomizedDelaySec=30

[Install]
WantedBy=timers.target
