# roles/pub/files/publish-morning-brief.service
[Unit]
Description=Sync and publish the Morning Brief (pub.lan)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/publish-morning-brief
